# config-stow/systemd-user/systemd/user/waybar.service.d/override.conf
# Drop-in over the distro's /usr/lib/systemd/user/waybar.service.
#
# Waybar is crash-prone around output add/remove (hotplug, DPMS, kanshi profile
# switches) — an upstream GTK bug, not a config error. The shipped unit only has
# Restart=on-failure and systemd's default start rate limit, so a burst of crashes
# during a single docking event leaves the bar dead until the next login.
#
# Restart=always also covers clean exits; disabling the rate limit keeps the bar
# coming back no matter how often it dies.
# Upstream crash: https://github.com/Alexays/Waybar/issues/3400

[Unit]
# StartLimit* are [Unit] options (systemd.unit(5)), not [Service].
# 0 disables start rate limiting entirely.
StartLimitIntervalSec=0

[Service]
Restart=always
RestartSec=1
