# config-stow/systemd-user/systemd/user/waybar.service.d/override.conf
# Drop-in over the distro's /usr/lib/systemd/user/waybar.service.
#
# Waybar is crash-prone around output add/remove (hotplug, DPMS, kanshi profile
# switches) — an upstream GTK bug, not a config error. The shipped unit only has
# Restart=on-failure and systemd's default start rate limit, so a burst of crashes
# during a single docking event leaves the bar dead until the next login.
#
# Restart=always also covers clean exits; disabling the rate limit keeps the bar
# coming back no matter how often it dies.
# Upstream crash: https://github.com/Alexays/Waybar/issues/3400

[Unit]
# StartLimit* are [Unit] options (systemd.unit(5)), not [Service].
# 0 disables start rate limiting entirely.
StartLimitIntervalSec=0

[Service]
Restart=always
RestartSec=1

# The #3400 GTK bug destroys the bar on ONE output *without killing the process*, so
# Restart=always alone never fires — the bar just goes half-dead until manually restarted
# (verified from the journal 2026-07-20: no process exit, no coredump, no logged assertion).
# Making GLib criticals fatal turns that `gtk_widget_set_visible: assertion 'GTK_IS_WIDGET'`
# (a g_return_if_fail critical) into an abort(): a coredump for the upstream issue, and a real
# process exit so Restart=always brings the whole bar back (~1s flicker, self-healing).
# Normal startup emits only a CSS *warning*, not a critical, so this does not crash-loop.
# NOTE (2026-07-21): catches ONLY the variant that emits a critical. The other observed variant
# — waybar logs a "successful" bar re-add after a `.configure` timeout that never renders — emits
# no critical, so this won't fire. General fix is restart-on-output-change; see
# Workstation-Documentation/desktop/waybar-supervision.md §3.
Environment=G_DEBUG=fatal-criticals
